# deploy/alphatrack-frontend.service
[Unit]
Description=AlphaTrack Frontend (Next.js)
After=network.target alphatrack-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/alphatrack
ExecStart=/usr/bin/node /opt/alphatrack/.next/standalone/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=PORT=4040
Environment=HOSTNAME=0.0.0.0

[Install]
WantedBy=multi-user.target
